# supervisor/supervisord.conf
; Run from repo root: supervisord -c supervisor/supervisord.conf

[unix_http_server]
file=/tmp/supervisor-parking.sock

[supervisord]
logfile=./logs/supervisord.log
logfile_maxbytes=5MB
logfile_backups=3
pidfile=/tmp/supervisord-parking.pid
childlogdir=./logs
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-parking.sock

[program:parking_app]
directory=%(ENV_PARKING_APP_DIR)s
command=%(ENV_PARKING_APP_PYTHON)s main.py
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
stdout_logfile=./logs/parking_app_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=./logs/parking_app_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1"
